# files/services/etc/systemd/user/db-secrets-init.service
# db-secrets-init — generate the DB-tier podman secrets for the rootless stack.
#
# This is a USER unit: it runs in the app user's systemd --user manager (started at
# boot via systemd-linger). It is enabled by a baked symlink in
# /etc/systemd/user/default.target.wants/ because `systemctl --user enable`
# cannot be run for the user at image build time.
#
# The postgres/valkey quadlets depend on this via Requires=/After= so the
# secrets exist before those containers start.

[Unit]
Description=Generate database-tier credentials (rootless podman secrets)
Documentation=https://github.com/nnnithinnn/alma-selfhosted
# No explicit ordering: the user manager has podman available by the time units
# start, and the DB quadlets pull this in via Requires=/After=db-secrets-init.
# (Avoid After=default.target here — combined with WantedBy=default.target on
# this unit and the dependent quadlets it would form an ordering cycle.)

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/db-secrets-init.sh

[Install]
WantedBy=default.target
