# files/system/etc/containers/systemd/tailscale.container
# Tailscale — node/subnet-router/exit-node client for this host's tailnet.
#
# STATUS: WIP / NOT YET ENABLED.
#   Only the image is pinned here so Renovate's `quadlet` manager tracks it.
#   No [Install] section -> generated but not started at boot.
#
# NOTE: running the Tailscale client in a container needs TUN access
#   (AddDevice=/dev/net/tun + NET_ADMIN) or userspace-networking mode, plus a
#   persistent state volume (/var/lib/tailscale). A host RPM daemon is the more
#   common pattern, but it is kept as a container here for uniform auto-bumping.
#
# Image is digest-pinned and auto-bumped by Renovate (.github/renovate.json).

[Unit]
Description=Tailscale node
Documentation=https://tailscale.com/

[Container]
Image=ghcr.io/tailscale/tailscale:v1.98.4@sha256:25cde9ad76020b0e29229136d0c38b5962e9a0e1774ffac9b0df68e4a37d6cf0
# TODO: AddDevice=/dev/net/tun, AddCapability=NET_ADMIN, state volume,
#       TS_AUTHKEY / TS_EXTRA_ARGS, etc.

[Service]
Restart=always

# TODO (enable later):
# [Install]
# WantedBy=multi-user.target default.target
